# dnsmasq/dnsmasq-supervisor.conf
[program:dnsmasq]
redirect_stderr=true
command=/usr/sbin/dnsmasq -C /etc/dnsmasq.conf --no-daemon
